# .config/autostart/i3lock.desktop
[Desktop Entry]
Exec=xautolock -locker i3lock
Name=i3lock
OnlyShowIn=LXQt;
Type=Application
Version=1.0
