# .config/autostart/i3lock.desktop
[Desktop Entry]
Exec=xautolock -locker \\\\\\"dm-tool lock\\\\\\"
Name=dm-tool lock
OnlyShowIn=LXQt;
Type=Application
Version=1.0
